# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; --------------------------------------------------
; PHP worker — restarts automatically if it ever dies
; --------------------------------------------------
[program:image-worker]
command=/usr/local/bin/php /app/worker.php
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; --------------------------------------------------
; Cron daemon — runs requeue-cron.php every minute
; --------------------------------------------------
[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.out.log